# scripts/com.smartir.autostream.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.smartir.autostream</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>bash</string>
        <string>-c</string>
        <string>
cd /Users/yu01/Desktop/SmartIR &amp;&amp; \
source .env &amp;&amp; \
/usr/local/bin/python3 scripts/run_auto_stream.py --type morning_market
        </string>
    </array>

    <!-- 平日朝 8:45 JST -->
    <key>StartCalendarInterval</key>
    <array>
        <!-- 月曜 -->
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>45</integer>
        </dict>
        <!-- 火曜 -->
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>45</integer>
        </dict>
        <!-- 水曜 -->
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>45</integer>
        </dict>
        <!-- 木曜 -->
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>45</integer>
        </dict>
        <!-- 金曜 -->
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>45</integer>
        </dict>
        <!-- 日曜 10:00 (週間まとめ) - 別のスクリプト実行が必要 -->
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/smartir-autostream.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/smartir-autostream-error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/yu01/Desktop/SmartIR</string>

    <!-- 失敗時のリトライを無効化 -->
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
<!--
インストール方法:
  cp scripts/com.smartir.autostream.plist ~/Library/LaunchAgents/
  launchctl load ~/Library/LaunchAgents/com.smartir.autostream.plist

アンインストール:
  launchctl unload ~/Library/LaunchAgents/com.smartir.autostream.plist
  rm ~/Library/LaunchAgents/com.smartir.autostream.plist

手動テスト:
  launchctl start com.smartir.autostream

ログ確認:
  tail -f /tmp/smartir-autostream.log
  tail -f /tmp/smartir-autostream-error.log

日曜日の週間まとめは別のplistを作成するか、
スクリプト内で曜日を判定して stream_type を切り替えてください。
-->
